# chezmoi/dot_config/systemd/user/arrpc.service
[Service]
Type=simple
ExecStart=/bin/bash -c '/usr/bin/node src'
WorkingDirectory=/usr/local/arrpc

[Install]
WantedBy=default.target